# branches/debug/etc/archie/rc.local
#!/bin/sh
#
# /etc/archie/rc.local: Local multi-user startup script.
#
#read KEY

# End of file
